# etc/systemd/system/nftables.service.d/load-firewall-sets.conf
[Service]
ExecStart=/etc/nftables/load-firewall-sets
ExecReload=/etc/nftables/load-firewall-sets
